# systemd/system/rsyslog.service.d/42-opt-in-allow-stretch-only.conf
# Workaround for Debian 9 / rsyslog 8.24.
# In Debian 10 / rsyslog 8.1901, -iNONE is supported and used by default.
[Service]
ReadWritePaths=/run
