# systemd_241_lockdown/etc/systemd/system/etckeeper.service.d/30-mitigate-service-denial.conf
# Upstream sets idle but not nice or batch.
[Service]
Nice=19
CPUSchedulingPolicy=batch
#IOSchedulingClass=idle
